# deploy/workspace/systemd/workspace-postgres-backup.service
[Unit]
Description=Back up the workspace PostgreSQL cluster
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=bukunmi
Group=docker
ExecStart=/opt/workspace/bin/backup-postgres
